# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/capwap.list
/usr/sbin/external_radius_login.sh
/sbin/check_wtp
/etc/capwap/client.pem
/etc/rc.d/rc.split_tunnel
/usr/sbin/setup_mtu.sh
/etc/capwap/check_data_channel.script
/usr/share/capwap/firewall.include
/etc/init.d/QosUserSrv
/sbin/WTP
/etc/init.d/capwap
/etc/capwap/settings.wtp.txt
/etc/capwap/root.pem
/etc/init.d/check_data_channel
/etc/uci-defaults/99-capwap
/etc/capwap/config.wtp
/usr/sbin/udptap_init.sh
/sbin/QosUserCli
/lib/WTP_module_EAP.so
/usr/sbin/split_tunnel_user.sh
/usr/sbin/get_ubus_sid.sh
/sbin/QosUserSrv
/etc/capwap/split_tun_action.script
/etc/rc.d/qos_user.sh
