# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/certificates.list
/etc/init.d/certificates
/lib/preinit/75_certificates
/etc/uci-defaults/99-tip-certificates.sh
