# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/cipuamd.list
/ramfs/bin/dec
/ramfs/uamd/cfg.proxy.logout.tmpl
/ramfs/uamd/proxy
/ramfs/uamd/external_login.tmpl
/ramfs/bin/handle_wizard.sh
/ramfs/uamd/punishweb.tmpl
/ramfs/uamd/cfg.proxy.fake.tmpl
/ramfs/bin/dhcpc.sh
/ramfs/bin/parseCN.sh
/ramfs/bin/enc
/ramfs/uamd/customize_page_data.tmpl
/ramfs/uamd/cfg.proxy.tmpl
/ramfs/bin/auth_customize_page.sh
/ramfs/uamd/redir_external_portal.tmpl
/ramfs/bin/run_setup_aaa.sh
/ramfs/bin/dns-nb
/usr/lib/libcipuamd.so
/etc/uci-defaults/98-cipuamd
/ramfs/bin/get_cert.sh
/ramfs/bin/wall_ipsetadt.sh
/ramfs/bin/run_uamd.sh
/ramfs/uamd/check_external_login.tmpl
/ramfs/uamd/rep_uamd_tmpl
/ramfs/uamd/microsoft365_login.tmpl
/etc/init.d/cipuamd
/usr/share/cipuamd/firewall.include
/ramfs/uamd/microsoft365_auth.tmpl
/ramfs/uamd/cfg.proxy.ticket.tmpl
