# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/client-limit.list
/usr/libexec/clientlimit.sh
/usr/bin/client-limit
/etc/hotplug.d/net/31-ratelimit-hs
/etc/hotplug.d/net/30-clientlimit
/usr/libexec/clientlimit-wait.sh
/etc/init.d/client-limit
