# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/conntrackd.list
/usr/sbin/conntrackd
/etc/conntrackd/conntrackd.conf
/lib/upgrade/keep.d/conntrackd
/etc/init.d/conntrackd
