# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/dhcp-helper.list
/usr/sbin/mkcfg-dhcp-helper.lua
/etc/init.d/dhcp_helper
/usr/sbin/dhcp-helper
