# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/discoveryd.list
/usr/sbin/discoveryd
/etc/init.d/discovery
