# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/dynamic-vlan.list
/usr/sbin/dynamic-vlan
/etc/init.d/dynamic-vlan
