# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/firstcontact.list
/usr/sbin/firstcontact
/etc/uci-defaults/zzz-firstcontact
/etc/init.d/firstcontact
/usr/share/ucentral/firstcontact.uc
