# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/hostapd-common.list
/etc/init.d/prc
/etc/capabilities/wpad.json
/usr/sbin/prc_relay.sh
/usr/share/acl.d/wpad_acl.json
/lib/netifd/hostapd.sh
/etc/init.d/wpad
/etc/rc.button/wps
/usr/sbin/ssdk_fdb_update.sh
