# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/luci-mod-status.list
/usr/share/luci/menu.d/luci-mod-status.json
/www/luci-static/resources/view/status/include/50_dsl.js
/www/luci-static/resources/view/status/connections.js
/usr/lib/lua/luci/view/admin_status/wirelessstat.htm
/www/luci-static/resources/view/status/include/20_memory.js
/www/luci-static/resources/wireless.svg
/www/luci-static/resources/view/status/include/genstat/20_system.js
/www/luci-static/resources/view/status/wireless.js
/sbin/logrotate_cron.sh
/www/luci-static/resources/view/status/bandwidth.js
/www/luci-static/resources/view/status/load.js
/www/luci-static/resources/view/status/include/40_dhcp.js
/usr/lib/lua/luci/view/admin_status/index.htm
/www/luci-static/resources/view/status/include/60_wifi.js
/www/luci-static/resources/view/status/include/netstat/30_routes.js
/usr/lib/lua/luci/controller/admin/acn_status.lua
/etc/init.d/send_clientinfo
/www/luci-static/resources/view/status/include/genstat/10_network.js
/www/luci-static/resources/view/status/index.js
/www/luci-static/resources/view/status/networkstat.js
/www/luci-static/resources/view/status/dmesg.js
/www/luci-static/resources/load.svg
/usr/lib/lua/luci/view/admin_status/wi_anchor.htm
/sbin/send_clientinfo_period.sh
/usr/lib/lua/luci/view/admin_status/syslog.htm
/usr/bin/luci-bwc
/sbin/get_binding_ssid.sh
/www/luci-static/resources/bandwidth.svg
/usr/lib/lua/luci/view/admin_status/generalstat.htm
/sbin/generate_wifi_client_info.sh
/www/luci-static/resources/view/status/syslog.js
/www/luci-static/resources/view/status/generalstat.js
/sbin/bwmonitor.sh
/www/luci-static/resources/wifirate.svg
/www/luci-static/resources/view/status/routes.js
/usr/lib/lua/luci/view/admin_status/networkstat.htm
/www/luci-static/resources/view/status/include/netstat/10_localnet.js
/www/luci-static/resources/view/status/include/10_system.js
/www/luci-static/resources/view/status/iptables.js
/www/luci-static/resources/view/status/processes.js
/www/luci-static/resources/view/status/include/wifistat/60_wifi.js
/www/luci-static/resources/connections.svg
/sbin/wi_anchor_curl.sh
/www/luci-static/resources/view/status/include/30_network.js
/www/luci-static/resources/view/status/include/netstat/20_dhcp.js
/www/luci-static/resources/view/status/wirelessstat.js
