# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/mac-auth.list
/usr/bin/mac_authd
/etc/init.d/mac-auth
/usr/sbin/login.sh
