# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/opennds.list
/usr/lib/opennds/get_client_token.sh
/etc/opennds/fas-aes.php
/usr/lib/opennds/unescape.sh
/usr/bin/ndsctl
/etc/opennds/fas-hid.php
/usr/lib/opennds/restart.sh
/etc/uci-defaults/40_opennds
/etc/opennds/htdocs/images/splash.jpg
/usr/lib/opennds/authmon.sh
/etc/opennds/fas-aes-https.php
/usr/lib/opennds/binauth_log.sh
/usr/lib/opennds/theme_click-to-continue-custom-placeholders.sh
/usr/bin/opennds
/etc/opennds/htdocs/css/captive-portal.css
/usr/lib/opennds/libopennds.sh
/etc/opennds/htdocs/css/bootstrap.min.css
/usr/lib/opennds/theme_user-email-login-basic.sh
/etc/config/opennds
/etc/opennds/htdocs/images/OpenWiFi_BadgeLogo_White.svg
/usr/lib/opennds/theme_click-to-continue-basic.sh
/usr/lib/opennds/get_client_interface.sh
/etc/opennds/htdocs/splash.css
/etc/opennds/htdocs/images/OpenWiFi_LogoLockup_DarkGreyColour.svg
/etc/init.d/opennds
/usr/lib/opennds/dnsconfig.sh
/usr/lib/opennds/theme_user-email-login-custom-placeholders.sh
/etc/opennds/config.uci
/usr/lib/opennds/client_params.sh
/usr/lib/opennds/post-request.php
