# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/phy-check.list
/usr/sbin/phy-check.sh
/etc/init.d/phy-check
