# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/qosify.list
/etc/qosify-defaults.conf
/etc/config/qosify
/lib/bpf/qosify-bpf.o
/etc/hotplug.d/net/10-qosify
/usr/sbin/qosify
/etc/init.d/qosify
