# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/syslog-ng.list
/usr/share/syslog-ng/include/scl/snmptrap/snmptrapd-source.conf
/usr/lib/syslog-ng/libafsocket.so
/usr/share/syslog-ng/include/scl/logscale/logscale.conf
/usr/share/syslog-ng/include/scl/sudo/sudo.conf
/usr/share/syslog-ng/include/scl/splunk/splunk.conf
/usr/bin/syslog-ng-update-virtualenv
/usr/lib/syslog-ng/libsecure-logging.so
/usr/share/syslog-ng/include/scl/cisco/plugin.conf
/usr/lib/syslog-ng/libafuser.so
/usr/lib/syslog-ng/libkvformat.so
/usr/lib/libsyslog-ng-4.4.so.0
/usr/lib/syslog-ng/libappmodel.so
/usr/share/syslog-ng/xsd/patterndb-3.xsd
/usr/lib/libloggen_plugin.so
/usr/share/syslog-ng/smart-multi-line.fsm
/usr/lib/syslog-ng/libbasicfuncs.so
/usr/lib/syslog-ng/libgraphite.so
/usr/lib/syslog-ng/libjson-plugin.so
/usr/lib/libloggen_helper.so
/usr/share/syslog-ng/include/scl/nodejs/plugin.conf
/usr/share/syslog-ng/include/scl/system/plugin.conf
/sbin/logread
/usr/sbin/syslog-ng-debun
/usr/lib/libloggen_helper-4.4.so.0
/usr/lib/libsecret-storage.so
/usr/lib/syslog-ng/libafprog.so
/usr/lib/syslog-ng/libtfgetent.so
/usr/share/syslog-ng/xsd/patterndb-4.xsd
/usr/share/syslog-ng/include/scl/slack/slack.conf
/usr/share/syslog-ng/include/scl/telegram/telegram.conf
/usr/share/syslog-ng/include/scl/kafka/kafka.conf
/usr/bin/pdbtool
/usr/lib/syslog-ng/libcryptofuncs.so
/usr/share/syslog-ng/include/scl/elasticsearch/elastic-java.conf
/usr/man/man1/syslog-ng-ctl.1
/usr/share/syslog-ng/include/scl/network_localhost/plugin.conf
/usr/share/syslog-ng/include/scl/apache/apache.conf
/usr/share/syslog-ng/include/scl/graphite/README
/usr/lib/syslog-ng/libmetrics-probe.so
/usr/lib/syslog-ng/libcorrelation.so
/usr/lib/syslog-ng/libaffile.so
/usr/bin/slogverify
/usr/share/syslog-ng/include/scl/collectd/plugin.conf
/usr/share/syslog-ng/xsd/patterndb-5.xsd
/usr/share/syslog-ng/include/scl/paloalto/panos.conf
/usr/share/syslog-ng/include/scl.conf
/usr/lib/libevtlog.so
/usr/man/man1/syslog-ng-debun.1
/usr/man/man5/syslog-ng.conf.5
/usr/share/syslog-ng/include/scl/kafka/kafka-java.conf
/usr/lib/syslog-ng/libxml.so
/usr/man/man1/pdbtool.1
/usr/lib/syslog-ng/libazure-auth-header.so
/usr/man/man1/slogencrypt.1
/usr/share/syslog-ng/include/scl/iptables/iptables.conf
/etc/syslog-ng.conf
/usr/share/syslog-ng/include/scl/junos/plugin.conf
/etc/init.d/syslog-ng
/usr/share/syslog-ng/tools/lex-rules.am
/usr/lib/syslog-ng/liblinux-kmsg-format.so
/usr/share/syslog-ng/include/scl/cee/adapter.conf
/usr/share/syslog-ng/include/scl/cim/adapter.conf
/usr/share/syslog-ng/xsd/patterndb-6.xsd
/usr/man/man1/loggen.1
/usr/lib/syslog-ng/libcef.so
/usr/lib/syslog-ng/libhook-commands.so
/usr/share/syslog-ng/include/scl/hdfs/plugin.conf
/usr/lib/libloggen_plugin-4.4.so.0.0.0
/usr/lib/libevtlog-4.4.so.0
/usr/share/syslog-ng/tools/cfg-grammar.y
/usr/share/syslog-ng/include/scl/osquery/plugin.conf
/usr/share/syslog-ng/include/scl/fortigate/fortigate.conf
/usr/share/syslog-ng/include/scl/solaris/plugin.conf
/usr/sbin/syslog-ng-ctl
/usr/share/syslog-ng/include/scl/graylog2/plugin.conf
/usr/share/syslog-ng/include/scl/syslogconf/plugin.conf
/usr/lib/syslog-ng/libtimestamp.so
/usr/share/syslog-ng/tools/merge-grammar.py
/usr/lib/libsyslog-ng.so
/usr/share/syslog-ng/include/scl/windowseventlog/plugin.conf
/usr/lib/syslog-ng/libconfgen.so
/usr/share/syslog-ng/include/scl/websense/plugin.conf
/usr/share/syslog-ng/include/scl/syslogconf/convert-syslogconf.awk
/usr/share/syslog-ng/include/scl/network_localhost/detect.sh
/usr/share/syslog-ng/include/scl/linux-audit/linux-audit.conf
/usr/share/syslog-ng/include/scl/checkpoint/plugin.conf
/usr/man/man8/syslog-ng.8
/usr/share/syslog-ng/include/scl/sumologic/sumologic.conf
/usr/man/man1/slogkey.1
/usr/lib/syslog-ng/loggen/libloggen_ssl_plugin.so
/usr/share/syslog-ng/include/scl/pacct/plugin.conf
/usr/share/syslog-ng/include/scl/python/python-modules.conf
/usr/lib/syslog-ng/libdisk-buffer.so
/usr/bin/update-patterndb
/usr/lib/syslog-ng/libadd-contextual-data.so
/usr/lib/syslog-ng/libregexp-parser.so
/usr/man/man1/slogverify.1
/usr/sbin/syslog-ng
/usr/lib/syslog-ng/libsystem-source.so
/usr/lib/syslog-ng/libsyslogformat.so
/usr/lib/libsyslog-ng-4.4.so.0.0.0
/usr/lib/libsecret-storage.so.0.0.0
/usr/lib/syslog-ng/libcsvparser.so
/etc/syslog-ng.d/.keep
/usr/share/syslog-ng/include/scl/opensearch/opensearch.conf
/usr/share/syslog-ng/include/scl/loggly/loggly.conf
/usr/bin/persist-tool
/usr/bin/loggen
/usr/share/syslog-ng/xsd/patterndb-1.xsd
/usr/lib/libsecret-storage.so.0
/usr/share/syslog-ng/include/scl/netskope/plugin.conf
/usr/man/man1/persist-tool.1
/usr/share/syslog-ng/include/scl/loadbalancer/plugin.conf
/usr/lib/libloggen_helper-4.4.so.0.0.0
/usr/lib/syslog-ng/librate-limit-filter.so
/usr/man/man7/secure-logging.7
/usr/lib/syslog-ng/libmap-value-pairs.so
/usr/lib/syslog-ng/libstardate.so
/usr/bin/dqtool
/usr/share/syslog-ng/include/scl/syslogconf/README
/usr/share/syslog-ng/include/scl/ewmm/ewmm.conf
/usr/share/syslog-ng/include/scl/mbox/mbox.conf
/lib/upgrade/keep.d/syslog-ng
/usr/bin/slogkey
/usr/share/syslog-ng/include/scl/loadbalancer/gen-loadbalancer.sh
/usr/share/syslog-ng/include/scl/elasticsearch/elastic-http.conf
/usr/lib/libloggen_plugin-4.4.so.0
/usr/lib/syslog-ng/libafstomp.so
/usr/lib/libevtlog-4.4.so.0.0.0
/usr/share/syslog-ng/include/scl/logmatic/logmatic.conf
/usr/share/syslog-ng/include/scl/rewrite/cc-mask.conf
/usr/lib/syslog-ng/libexamples.so
/usr/share/syslog-ng/include/scl/windowseventlog/windowseventlog.xml
/usr/share/syslog-ng/xsd/patterndb-2.xsd
/usr/man/man1/dqtool.1
/usr/share/syslog-ng/include/scl/graphite/plugin.conf
/usr/lib/syslog-ng/libhttp.so
/usr/lib/syslog-ng/libpseudofile.so
/usr/share/syslog-ng/tools/system-expand.sh
/usr/share/syslog-ng/include/scl/default-network-drivers/plugin.conf
/usr/bin/slogencrypt
/usr/share/syslog-ng/include/scl/mariadb/audit.conf
/usr/lib/syslog-ng/libtags-parser.so
/usr/share/syslog-ng/include/scl/discord/discord.conf
/usr/share/syslog-ng/include/scl/cim/template.conf
/usr/lib/syslog-ng/loggen/libloggen_socket_plugin.so
