# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/tisbl.list
/usr/bin/com-wr.sh
/usr/bin/tisbl
/usr/bin/tisbl.sh
/usr/bin/ble-scan-rx-parser.sh
/usr/bin/hexdump.sh
/etc/tifirmware/multi_role_bd9.bin
/etc/tifirmware/znp_bd9.bin
/etc/tifirmware/ble5_host_test_sdk6_bd9.bin
/etc/config/tisbl
/usr/bin/blescan.sh
/etc/tifirmware/blinky_bd9.bin
/etc/uci-defaults/118-tisbl.defaults
/usr/bin/ibeacon-scan.sh
/usr/bin/ibeacon-broadcaster.sh
/etc/tifirmware/simple_broadcaster_bd9.bin
/etc/config/ibeacon
/usr/bin/zbcli.sh
/etc/init.d/tisbl
/etc/tifirmware/simple_central_bd9.bin
/etc/tifirmware/simple_peripheral_bd9.bin
/usr/bin/blescan2file.sh
/etc/init.d/ibeacon
/etc/uci-defaults/119-ibeacon.defaults
/usr/bin/wdt_touch.sh
/usr/bin/ti-hosttest2net.sh
/etc/tifirmware/ble5_host_test_bd9.bin
