# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/ucentral-event.list
/etc/config/event
/usr/sbin/ucentral-event
/etc/init.d/ucentral-event
