# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/ucentral-wifi.list
/usr/sbin/ucentral-wifi
/etc/init.d/ucentral-wifi
