# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/ulog-listener.list
/etc/init.d/ulog-listener
/sbin/ulog-listener
