# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/umdns.list
/etc/config/umdns
/etc/uci-defaults/58_umdns-add-wan-support
/usr/sbin/umdns
/etc/init.d/umdns
/etc/seccomp/umdns.json
