# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/usr/lib/opkg/info/wifischedule.list
/etc/rc.d/S43wifischedule
/etc/config/wifi_schedule
/etc/init.d/wifischedule
/usr/bin/wifi_schedule.sh
