# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/etc/uci-defaults/58_umdns-add-wan-support
#!/bin/sh

. /lib/functions.sh


if [ "$(uci -q get umdns.@umdns[-1].network | grep wan)" == "" ]; then
	uci -q batch <<-EOF >/dev/null
	add_list umdns.@umdns[-1].network=wan
	EOF
	uci_commit umdns
fi

exit 0
