# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/etc/uci-defaults/61-wizard-acn
#!/bin/sh

if [ "$(uci -q get acn.wizard)" != "ui" ]; then
uci -q batch <<-EOF >/dev/null
  set acn.wizard=ui
  set acn.wizard.enabled='1'
  commit acn.wizard
EOF
fi

exit 0
