# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/etc/uci-defaults/80-register-port-default
#!/bin/sh

if [ "$(uci -q get acn.register.port)" == "" ]; then
	uci set acn.register.port='5222,443'
	uci commit acn.register
fi
exit 0