# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/init.d/acc
#!/bin/sh /etc/rc.common

START=37

start() {
	/usr/sbin/acc hw all >/dev/null
}
