# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/init.d/check-warm-start
#!/bin/sh /etc/rc.common
STOP=90

USE_PROCD=1

WARM_START_FILE="/etc/warm_restart"

stop_service() {
    [ -f "$WARM_START_FILE" ] || touch "$WARM_START_FILE"
}
