# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
/sbin/channel_balance.sh boot
/sbin/wmi_monitor.sh boot
/sbin/logrotate_cron.sh boot
/sbin/wifi_monitor.sh boot

exit 0
