# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/uci-defaults/124-uhttp-default
#!/bin/sh

. /lib/functions.sh

uci set uhttpd.main.rfc1918_filter='0'
uci commit uhttpd.main

exit 0
