# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/uci-defaults/53_network-inet-src
#!/bin/sh

if [ "$(uci -q get network.wan.inet_src)" == "" ]; then
	uci set network.wan.inet_src='eth0'
	uci commit network.wan
fi

exit 0
