# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/uci-defaults/65-dhcpc-hostname
#!/bin/sh

if [ "$(uci -q get network.wan.hostname)" == "" ]; then
	uci set network.wan.manual_hostname='1'
	uci set network.wan.hostname='Edge-core'
        uci commit network.wan
fi

exit 0
