# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/uci-defaults/77-igmp-snooping-default
#!/bin/sh

if [ "$(uci -q get network.lan)" == "" ]; then
	uci set network.lan=interface
fi

uci set network.lan.igmp_snooping='0'

uci commit network.lan

exit 0