# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/uci-defaults/77-msp-enable
#!/bin/sh

if [ "$(uci -q get acn.settings.msp_enable)" == "" ]; then
        uci set acn.settings.msp_enable='0'
fi

exit 0