# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/uci-defaults/78-discovery-tool-default
#!/bin/sh

if [ "$(uci -q get acn.discovery)" == "" ]; then
        uci set acn.discovery='acn'
        uci set acn.discovery.enabled='1'
fi

uci commit acn.discovery

exit 0