# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/uci-defaults/99-ucentral-umdns
#!/bin/sh

tip_mode=$(uci -q get acn.mgmt.management)
[ ! -z $tip_mode ] && [ "$tip_mode" -eq 3 ] || {
	exit 0
}

uci delete umdns.@umdns[0].network
