# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap104/etc/uci-defaults/luci-theme-ec
#!/bin/sh
uci batch <<-EOF
	set luci.themes.ec=/luci-static/ec
	set luci.main.mediaurlbase=/luci-static/ec
	commit luci
EOF
exit 0
