# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/acc.list
/usr/sbin/trouble.sh
/usr/sbin/acc
/etc/init.d/acc
