# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/cipfwd.list
/sbin/cipfwd
/etc/firewall.dai
/sbin/sender
/etc/uci-defaults/57_firewall-dai
/etc/init.d/cipfwd
