# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/cipgwacct.list
/ramfs/bin/hostapd_8021x.sh
/ramfs/bin/cipgwcli
/ramfs/bin/cipgwnlsrv
/ramfs/bin/cipgwNewIP.sh
/ramfs/bin/cipgwsrv
/usr/lib/libcipgwacct.so
/etc/init.d/cipgwnlsrv
/etc/init.d/cipgwsrv
