# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/ec-fp.list
/etc/hotplug.d/net/10-ecfp
/usr/sbin/ec-fp
/etc/init.d/ec-fp
/usr/share/ecfp/dhcp.json
