# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/ieee8021x.list
/usr/sbin/ieee8021x
/etc/config/ieee8021x
/etc/init.d/ieee8021x
