# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/kmod-qca-ssdk.list
/etc/init.d/qca-ssdk
/lib/modules/4.4.60/qca-ssdk.ko
/etc/modules.d/30-qca-ssdk
