# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/luci-mod-system.list
/etc/uci-defaults/112-allow-over-wan-ec-discovery-tool
/etc/uci-defaults/74-telnetd
/usr/lib/lua/luci/model/cbi/user.lua
/etc/uci-defaults/111-fix-eap104-wifi-led
/etc/uci-defaults/77-igmp-snooping-default
/www/luci-static/resources/view/system/flash.js
/usr/lib/lua/luci/model/cbi/services.lua
/usr/lib/lua/luci/view/admin_system/easy_wizard.htm
/www/luci-static/resources/view/system/sshkeys.js
/usr/lib/lua/luci/view/admin_system/wizard.htm
/www/luci-static/resources/view/system/password.js
/usr/lib/lua/luci/model/cbi/admin_system/backupfiles.lua
/etc/uci-defaults/64-wizard-qrcode
/etc/uci-defaults/66_leds
/lib/acn/acn_functions.sh
/www/luci-static/resources/view/system/system.js
/www/luci-static/resources/view/system/dropbear.js
/etc/uci-defaults/63-default-country
/usr/lib/lua/luci/view/admin_system/upgrade.htm
/www/luci-static/resources/view/system/crontab.js
/etc/init.d/users
/lib/preinit/83_104_mac_check
/etc/uci-defaults/70-remote-syslog-default
/sbin/reset_to_defaults.sh
/usr/sbin/config_restore.sh
/lib/acn/post-restore.sh
/usr/lib/lua/luci/controller/admin/system.lua
/etc/uci-defaults/121-system-time-zone
/www/luci-static/resources/view/system/mounts.js
/etc/uci-defaults/77-reset-and-led-button
/etc/uci-defaults/116_eap101_eap102_leds
/etc/uci-defaults/120-eap104_low_cost
/www/luci-static/resources/view/system/leds.js
/etc/uci-defaults/65-dhcpc-hostname
/sbin/enable_network_trigger.sh
/www/luci-static/resources/view/system/service.js
/etc/uci-defaults/69-lldp-default
/www/luci-static/resources/view/system/reboot.js
/usr/lib/lua/luci/view/admin_system/certificate.htm
/etc/uci-defaults/77-msp-enable
/usr/lib/lua/luci/view/admin_system/applyreboot.htm
/sbin/rpcd_add_accout.sh
/sbin/ec_firstcontact.sh
/etc/rc.d/S91users
/etc/hotplug.d/iface/90-ec-firstcontact
/etc/uci-defaults/76-cloudsite-India
/lib/preinit/82_ec_check_srv
/etc/uci-defaults/80-register-port-default
/etc/uci-defaults/75-mgmt-migrate
/etc/uci-defaults/79-led-power-default
/etc/config/hotspot20
/etc/uci-defaults/125-syslog-level
/usr/lib/lua/luci/model/cbi/system.lua
/etc/uci-defaults/124-uhttp-default
/etc/uci-defaults/78-discovery-tool-default
/etc/uci-defaults/68-capwap-default
/www/luci-static/resources/view/system/startup.js
/usr/lib/lua/luci/view/admin_system/qr_wizard.htm
/etc/uci-defaults/62-users-default
/usr/share/luci/menu.d/luci-mod-system.json
/etc/uci-defaults/61-wizard-acn
