# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/maverick.list
/etc/init.d/maverick
/usr/libexec/ucentral/maverick.sh
/usr/sbin/maverick
