# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/mesh-discovery.list
/usr/sbin/mesh-discovery
/usr/sbin/sync_wireless_vaps.sh
/etc/init.d/meshd
