# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/ser2net.list
/etc/ser2net.conf
/etc/config/ser2net
/etc/init.d/ser2net
/usr/sbin/ser2net
