# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/service-schedule.list
/usr/sbin/service_schedule.sh
/etc/init.d/service-schedule
