# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/usr/lib/opkg/info/udhcpsnoop.list
/etc/config/dhcpsnoop
/etc/init.d/dhcpsnoop
/usr/sbin/udhcpsnoop
