# EAP101_SRC/EDGECORE_AP-SRC/sysupgrade-edgecore_eap101/etc/hotplug.d/iface/90-ec-firstcontact
#!/bin/sh

[ "$ACTION" == "ifup" ] || exit 0

[ "$INTERFACE" = "wan" ] && {
	/sbin/ec_firstcontact.sh
}
